# templates/healthchecker.service.j2
[Unit]
Description=PowerDNS Health Checker Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/healthchecker/src
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=healthchecker
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target