# files/dns-healthcheck.service
[Unit]
Description=DNS Health Check Service
After=network.target

[Service]
Type=simple
User=healthcheck
Group=healthcheck
WorkingDirectory=/opt/healthchecker/src
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dns-healthcheck
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target